# packages/cli/launchd/com.sync-tasks.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Template file - DO NOT USE DIRECTLY
  Run scripts/install.sh to generate the actual plist with your configuration
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sync-tasks</string>

    <key>ProgramArguments</key>
    <array>
        <string>/path/to/sync-tasks</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>WEBHOOK_URL</key>
        <string>YOUR_WEBHOOK_URL</string>
        <key>WEBHOOK_SECRET</key>
        <string>YOUR_SECRET_HERE</string>
    </dict>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/sync-tasks.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/sync-tasks.error.log</string>
</dict>
</plist>
